# docker/supervisor/vite.conf
[program:vite]
command=npm run dev
directory=/var/www
user=root
autostart=true
autorestart=true